# misp_feeds/The_Botvrij.eu_Data/5fd62412-4cec-450a-a2bc-e52ec0a8ab16.json
{"Event": {"timestamp": "1607869658", "threat_level_id": "2", "published": true, "date": "2020-12-13", "uuid": "5fd62412-4cec-450a-a2bc-e52ec0a8ab16", "publish_timestamp": "1611866658", "analysis": "2", "info": "xHunt Campaign: Newly Discovered Backdoors Using Deleted Email Drafts and DNS Tunneling for Command and Control", "Orgc": {"uuid": "56c42374-fdb8-4544-a218-41ffc0a8ab16", "name": "CUDESO"}, "Tag": [{"colour": "#00d622", "exportable": true, "name": "tlp:white"}, {"colour": "#0088cc", "exportable": true, "name": "misp-galaxy:mitre-enterprise-attack-attack-pattern=\"Scheduled Task - T1053\""}], "Attribute": [{"timestamp": "1607869509", "comment": "TriFive Samples", "category": "Payload delivery", "value": "407e5fe4f6977dd27bc0050b2ee8f04b398e9bd28edd9d4604b782a945f8120f", "type": "sha256", "uuid": "5fd62445-f8a4-4cfc-ba35-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy Samples", "category": "Payload delivery", "value": "c18985a949cada3b41919c2da274e0ffa6e2c8c9fb45bade55c1e3b6ee9e1393", "type": "sha256", "uuid": "5fd62445-4a78-442d-b8c1-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy Samples", "category": "Payload delivery", "value": "6c13084f213416089beec7d49f0ef40fea3d28207047385dda4599517b56e127", "type": "sha256", "uuid": "5fd62445-7840-4955-b15b-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy Samples", "category": "Payload delivery", "value": "efaa5a87afbb18fc63dbf4527ca34b6d376f14414aa1e7eb962485c45bf38372", "type": "sha256", "uuid": "5fd62445-afcc-434c-8726-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy Samples", "category": "Payload delivery", "value": "a4a0ec94dd681c030d66e879ff475ca76668acc46545bbaff49b20e17683f99c", "type": "sha256", "uuid": "5fd62445-2830-458f-91b6-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy C2 Domains", "category": "Network activity", "value": "deman1.icu", "type": "domain", "uuid": "5fd62445-f940-415f-9112-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy C2 Domains", "category": "Network activity", "value": "hotsoft.icu", "type": "domain", "uuid": "5fd62445-d928-45b3-9e65-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy C2 Domains", "category": "Network activity", "value": "uplearn.top", "type": "domain", "uuid": "5fd62445-c938-44b3-a222-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy C2 Domains", "category": "Network activity", "value": "lidarcc.icu", "type": "domain", "uuid": "5fd62445-6678-4132-bdc5-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869509", "comment": "Snugy C2 Domains", "category": "Network activity", "value": "sharepoint-web.com", "type": "domain", "uuid": "5fd62445-a250-417b-9807-e42ec0a8ab16", "to_ids": true, "object_relation": null}, {"timestamp": "1607869524", "comment": "", "category": "External analysis", "value": "https://unit42.paloaltonetworks.com/xhunt-campaign-backdoors/", "type": "url", "uuid": "5fd62454-5e70-452b-8f16-e42fc0a8ab16", "to_ids": false, "object_relation": null}]}}